# root/etc/init.d/jianguoyun-backup
#!/bin/sh /etc/rc.common
# 坚果云备份插件 - init.d服务脚本
# 用于管理定时任务，确保重启后备份计划保留

START=99
STOP=10

PROG=/usr/bin/jianguoyun-backup.sh

start_service() {
    # 服务启动时设置定时任务
    if [ -x "$PROG" ]; then
        $PROG setup_cron >/dev/null 2>&1
    fi
}

stop_service() {
    # 服务停止时移除定时任务（使用 crontab 命令，更规范）
    crontab -l 2>/dev/null | grep -v 'jianguoyun-backup' | crontab - 2>/dev/null
    /etc/init.d/cron restart 2>/dev/null
}

restart_service() {
    stop_service
    start_service
}

reload_service() {
    # 重新加载配置（重新设置定时任务）
    if [ -x "$PROG" ]; then
        $PROG setup_cron >/dev/null 2>&1
    fi
}

service_running() {
    # 检查定时任务是否存在
    crontab -l 2>/dev/null | grep -q 'jianguoyun-backup'
}

boot() {
    # 开机时启动服务
    start_service
}
